# xray-bot.service
[Unit]
Description=Xray Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/xray-bot
Environment="PATH=/opt/xray-bot/venv/bin"
EnvironmentFile=/opt/xray-bot/.env
ExecStart=/opt/xray-bot/venv/bin/python /opt/xray-bot/bot.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target